# isyhelper.service
[Unit]
Description=ISYHelper Server
Documentation=https://github.com/jimboca/ISYHelper/blob/master/README.md
After=network-online.target

[Service]
Type=simple
<<<<<<< HEAD
WorkingDirectory=/home/pi/isyhelper/ISYHelper
=======
ExecStartPre=cd /home/pi/isyhelper/ISYHelper
>>>>>>> dev
ExecStart=/home/pi/isyhelper/ISYHelper/isyhelper.py
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
